# multi-arch/packages/srv-k3s/services/k3s/run
#!/bin/sh

exec 2>&1
img="/k3s.img"

if [ ! -d "/var/lib/rancher" ]; then
 mkdir -p /var/lib/rancher
 if [ ! -e "$img" ]; then
   dd if=/dev/zero of=$img bs=1MB count=2048
   mkfs.ext4 $img
   tune2fs -c0 -i0 $img
 fi
 mount -t ext4 -o loop "$img" /var/lib/rancher
fi

echo 'k3s starting.'

LANG=C LC_ALL=C \
exec /usr/bin/k3s server 2>&1
