# multi-arch/packages/yip-systemd/yip-reconcile.timer
[Unit]
Description=Yip reconcile
 
[Timer]
OnBootSec=5min
OnUnitActiveSec=60min
Unit=yip-reconcile.service
 
[Install]
WantedBy=multi-user.target
