# multi-arch/packages/yip-systemd/yip-boot.service
[Unit]
Description=Yip system setup
Before=display-manager.service getty.target

[Service]
Type=oneshot
#RemainAfterExit=true
ExecStart=/bin/bash -c "systemd-inhibit /usr/bin/yip-init boot"

[Install]
WantedBy=multi-user.target