# packages/srv-udhcpc/services/udhcpc/run
#!/bin/sh

exec 2>&1

echo 'udhcpc starting.'

LANG=C LC_ALL=C \
exec /sbin/udhcpc -f -i "${INTERFACE}" 2>&1
