# packages/srv-k3s/services/k3s/run
#!/bin/sh

exec 2>&1

echo 'k3s starting.'

LANG=C LC_ALL=C \
exec /usr/bin/k3s server 2>&1
